# 2020-10-25_defcon-b3s23/systemd/challenge.socket
[Unit]
Description=Challenge Socket

[Socket]
ListenStream=2323
#NoDelay=true
#FreeBind=true
Accept=true

[Install]
WantedBy=sockets.target
